# usb_autorun.service
[Unit]
Description=USB Autorun Service for the Memory Device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb_autorun.sh
